# systemd/orgblog.service
[Unit]
Description=Orgblog API

[Service]
ExecStart=/usr/local/bin/orgblog

[Install]
WantedBy=multi-user.target